# resources/installer-scripts/container/common/container-hal9000-console-proxy.socket.template
[Unit]
Description=container-hal9000-console-proxy.socket
BindsTo=pod-hal9000.service

[Socket]
ListenStream=SYSTEM_IP:8080

[Install]
WantedBy=sockets.target

